# systemd/openclaw-proxy.service
[Unit]
Description=OpenClaw Keymaster Proxy
Documentation=https://github.com/dommurphy155/Keymaster
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%I
WorkingDirectory=/home/%I/.openclaw/skills/keymaster
Environment=PYTHONPATH=/home/%I/.openclaw/skills/keymaster
Environment=HOME=/home/%I
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/%I/.local/bin

# Pre-start: Configure OpenClaw to use proxy
ExecStartPre=/bin/bash -c 'echo "[$(date)] Configuring OpenClaw to use proxy..." >> /home/%I/.openclaw/keymaster_proxy.log'
ExecStartPre=/usr/bin/python3 /home/%I/.openclaw/skills/keymaster/scripts/configure_openclaw.py --enable

# Start the proxy with unbuffered output for better logging
ExecStart=/usr/bin/python3 -u -m uvicorn proxy.server:app --host 127.0.0.1 --port 8787 --loop uvloop --log-level info

# Log output
StandardOutput=append:/home/%I/.openclaw/keymaster_proxy.log
StandardError=append:/home/%I/.openclaw/keymaster_proxy.log

# Restart on failure
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Graceful shutdown
TimeoutStopSec=30
KillSignal=SIGTERM

# File descriptor and process limits
LimitNOFILE=65535
LimitNPROC=4096
TasksMax=infinity

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/%I/.openclaw

[Install]
WantedBy=multi-user.target
